# qinglong/rootfs/etc/s6-overlay/s6-rc.d/qinglong/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Qinglong 启动脚本 (S6 longrun service)
# 在 cont-init.d 完成持久化目录设置后，启动青龙面板
# ==============================================================================

bashio::log.info "Starting Qinglong..."

export QL_DATA_DIR=/config/qinglong/data
cd /ql
exec ./docker/docker-entrypoint.sh
