# lucky/rootfs/etc/s6-overlay/s6-rc.d/lucky/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Lucky 启动脚本 (S6 longrun)
# 在 cont-init.d 完成持久化目录设置后，启动 Lucky（前台阻塞运行）
# ==============================================================================
bashio::log.info "Starting Lucky..."
exec /usr/local/bin/lucky -c /config/luckyconf/lucky.conf
