# lucky/rootfs/etc/s6-overlay/s6-rc.d/lucky/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Lucky 启动脚本 (S6 longrun)
# 在 cont-init.d 完成持久化目录设置后，启动 Lucky（前台阻塞运行）
# ==============================================================================
bashio::log.info "Starting Lucky..."

/usr/local/bin/lucky -c /config/luckyconf/lucky.conf &
APP_PID=$!
trap 'kill $APP_PID 2>/dev/null; exit 0' TERM INT
bashio::log.info "Waiting for Lucky to become ready..."
for i in $(seq 1 60); do
    if curl -sf http://127.0.0.1:16601/ 2>/dev/null; then
        bashio::log.info "Lucky is ready (${i}s)"
        printf '\n' >&3
        wait $APP_PID
        exit $?
    fi
    sleep 1
done
bashio::log.error "Lucky failed to start within 60s"
kill $APP_PID 2>/dev/null
exit 1
