# baihu-panel/rootfs/etc/s6-overlay/s6-rc.d/baihu/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Baihu Panel 启动脚本 (S6 longrun service)
# 在 cont-init.d 完成持久化目录设置后，启动白虎面板
# ==============================================================================

bashio::log.info "Starting Baihu Panel..."
exec /usr/local/bin/baihu
