# qdtoday/rootfs/etc/s6-overlay/s6-rc.d/qdtoday/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# QD-Today 启动脚本 (S6 longrun)
# 在 cont-init.d 完成持久化目录设置后，启动 QD-Today Tornado 应用（前台阻塞）
# ==============================================================================
bashio::log.info "Starting QD-Today..."

# 仅在用户显式设置配置项时导出对应环境变量，留空则由 QD 使用内置默认值
_set_if_configured() {
  local key="$1"
  if bashio::config.has_value "${key}"; then
    local val
    val="$(bashio::config "${key}")"
    export "${key}=${val}"
    bashio::log.debug "Configured: ${key}=${val}"
  fi
}

_set_if_configured "QD_DEBUG"
_set_if_configured "AES_KEY"
_set_if_configured "COOKIE_SECRET"
_set_if_configured "DOMAIN"
_set_if_configured "USE_PYCURL"
_set_if_configured "DNS_SERVER"
_set_if_configured "TRACEBACK_PRINT"
_set_if_configured "PUSH_BATCH_SW"
_set_if_configured "PROXIES"
_set_if_configured "MAIL_SMTP"
_set_if_configured "MAIL_PORT"
_set_if_configured "MAIL_SSL"
_set_if_configured "MAIL_USER"
_set_if_configured "MAIL_PASSWORD"
_set_if_configured "TASK_REQUEST_LIMIT"
_set_if_configured "TASK_WHILE_LOOP_TIMEOUT"
_set_if_configured "NEW_TASK_DELAY"

cd /usr/src/app
exec python3 run.py
