# filebrowser-quantum/rootfs/etc/s6-overlay/s6-rc.d/filebrowser-quantum/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# FileBrowser Quantum 启动脚本 (S6 longrun)
# ==============================================================================
bashio::log.info "Starting FileBrowser Quantum..."
exec /usr/local/bin/filebrowser --config /config/filebrowser-quantum/config.yaml